# examples/systemd-with-container/owncast-rtmp@unixdemo.socket.d/override.conf
[Socket]
ListenStream=%h/owncast-rtmp-socket.%i
