# examples/systemd-with-container/owncast-webserver@unixdemo.socket.d/override.conf
[Socket]
ListenStream=%h/owncast-webserver-socket.%i
